# recipes-appends/oe-core/watchdog/watchdog.initd
#!/sbin/openrc-run

depend() {
	need localmount
	use logger
}

command_args="${COMMAND_ARGS}"
command_args_foreground="-F"
command="/usr/sbin/watchdog"
supervisor=supervise-daemon

# vim: noet ft=gentoo-init-d
