# recipes-init/openrc/openrc/readonly-rootfs.initd
#!/sbin/openrc-run
description="read-only rootfs remount"

depend() {
	before *
}

start() {
	mount -o remount,ro /
}
