# recipes-init/openrc-inittab-busybox/openrc-inittab-busybox/inittab.in
::sysinit:${sbindir}/openrc sysinit
::sysinit:${sbindir}/openrc boot
::wait:${sbindir}/openrc default
::shutdown:${sbindir}/openrc shutdown
::ctrlaltdel:${sbindir}/openrc reboot
